# dotfiles_linux/systemd/.config/systemd/user/waybar.service.d/workspace-activity.conf
# Generate the workspace-activity CSS before waybar starts, so the @import
# in style.css always has a file to load.
[Service]
ExecStartPre=%h/bin/workspace-activity render
