# dotfiles_linux/systemd/.config/systemd/user/waybar.service.d/workspace-colors.conf
# Generate the per-repo workspace-colors CSS before waybar starts, so the
# @import in style.css always has a file to load.
[Service]
ExecStartPre=%h/bin/waybar-workspace-colors
