# dotfiles_linux/systemd/.config/systemd/user/uair.service
[Unit]
Description=Uair

[Service]
Type=simple
ExecStart=/usr/bin/uair -q

[Install]
WantedBy=default.target

